# build/package/deb/postinst
#!/bin/bash

PACKAGE_NAME="vagrant-private-cloud"
SERVICE_NAME="vagrant-private-cloud"
BASE_DIR="/opt/vagrant-private-cloud"

function log() {
    echo -e "[$PACKAGE_NAME][${USER}][`date`] - ${*}"
}

set -e

OP="$1"

log "[postinst]" Running "$OP" for package "$PACKAGE_NAME"

case "${OP}" in
    install|configure)
    
        if [ ! -d ${BASE_DIR}/storage ]; then
            mkdir -p ${BASE_DIR}/storage
        fi
        adduser --system --no-create-home --group vagrant || true
        usermod -a -G syslog vagrant 
        chown vagrant:vagrant ${BASE_DIR} -R
        systemctl enable "${SERVICE_NAME}" >/dev/null || true
        systemctl start "${SERVICE_NAME}"|| log "Failed to start ${SERVICE_NAME} service"

    ;;
    update|upgrade)

        if [ ! -d ${BASE_DIR}/storage ]; then
            mkdir -p ${BASE_DIR}/storage
        fi
        chown vagrant:vagrant ${BASE_DIR} -R
        systemctl --system daemon-reload > /dev/null || true
        systemctl restart "${SERVICE_NAME}" || log "Failed to restart ${SERVICE_NAME} service"

    ;;
    abort-upgrade)
        # do nothing
    ;;
    *)
        log "postinst called with unknown argument" ${OP} >&2
        exit 1
    ;;
esac

exit 0

